# systemd/mta-yoklama.service
[Unit]
Description=MTA Yoklama Next.js Uygulamasi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mtaadmin
Group=mtaadmin
WorkingDirectory=/home/mtaadmin/Desktop/MtaYoklama
Environment=NODE_ENV=production
Environment=SESSION_COOKIE_SECURE=false
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target